# src/mymcp/deploy/templates/mymcp.service.in
[Unit]
Description=MyMCP Server
After=network.target

[Service]
Type=simple
User={service_user}
WorkingDirectory={working_directory}
EnvironmentFile={env_file}
ExecStart={exec_start}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
