# deploy/mymcp.service
[Unit]
Description=MyMCP Linux Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mymcp
EnvironmentFile=/opt/mymcp/.env
ExecStart=/opt/mymcp/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8765
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
